# linux/systemd/FOSSLinux.service
[Unit]
Description=FOSSLinux Service

[Service]
Type=notify
ExecStart=/usr/local/bin/Hello_FOSSLinux.sh

[Install]
WantedBy=multi-user.target
